# systemd/weekly-discord-bot.service
[Unit]
Description=Weekly Transaction Pipeline Discord Bot
After=network.target

[Service]
Type=simple
User=asya
WorkingDirectory=/mnt/DATA/Proyek/agency/discord-bot-weekly
ExecStart=/usr/bin/env node index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weekly-discord-bot

# Loads environment variables from the .env file
EnvironmentFile=/mnt/DATA/Proyek/agency/discord-bot-weekly/.env

[Install]
WantedBy=multi-user.target
